# packaging/launchd/tech.hyperbliss.hypercolor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>tech.hyperbliss.hypercolor</string>

    <key>ProgramArguments</key>
    <array>
        <string>@BIN_DIR@/hypercolor-daemon</string>
        <string>--ui-dir</string>
        <string>@UI_DIR@</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>StandardOutPath</key>
    <string>@LOG_DIR@/hypercolor.log</string>

    <key>StandardErrorPath</key>
    <string>@LOG_DIR@/hypercolor.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HYPERCOLOR_LOG</key>
        <string>info</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:@BIN_DIR@</string>
    </dict>

    <key>ProcessType</key>
    <string>Standard</string>

    <key>LowPriorityBackgroundIO</key>
    <true/>
</dict>
</plist>
